# deploy/systemd/tentacle-nftables.service
[Unit]
Description=Tentacle Firewall Manager
Documentation=https://github.com/joyautomation/tentacle
After=tentacle-nats.service
Requires=tentacle-nats.service

[Service]
Type=simple
EnvironmentFile=/opt/tentacle/config/tentacle.env
WorkingDirectory=/opt/tentacle/services/tentacle-nftables
ExecStart=/opt/tentacle/bin/deno run --allow-all main.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tentacle-nftables

[Install]
WantedBy=multi-user.target
