# deploy/systemd/tentacle-mqtt.service
[Unit]
Description=Tentacle MQTT Sparkplug B Bridge
Documentation=https://github.com/joyautomation/tentacle
After=tentacle-nats.service
Requires=tentacle-nats.service

[Service]
Type=simple
EnvironmentFile=/opt/tentacle/config/tentacle.env
WorkingDirectory=/opt/tentacle/services/tentacle-mqtt
ExecStart=/opt/tentacle/bin/deno run --allow-all main.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tentacle-mqtt

[Install]
WantedBy=multi-user.target
